# deploy/vpn/linkedin-mcp.service.d/vpn.conf
# Drop-in for linkedin-mcp.service that binds the MCP to the split-tunnel.
#
# Install at /etc/systemd/system/linkedin-mcp.service.d/vpn.conf, then
# `systemctl daemon-reload && systemctl restart linkedin-mcp`.
#
# Requires+After pull the tunnel up first; the Environment line points
# Chromium's proxy at the in-namespace tinyproxy so all LinkedIn traffic
# egresses through the university IP.
[Unit]
Requires=linkedin-vpn.service
After=linkedin-vpn.service

[Service]
Environment=LINKEDIN_PROXY_SERVER=http://10.200.0.2:8888
# Enable the VM-only VPN self-heal MCP tools (vpn_status / vpn_reconnect /
# vpn_egress_ip). Omit to run without them.
Environment=LINKEDIN_VPN_ENABLED=1
